# Custom Detection-as-Code Pipeline (Sigma + SIEM)/rules/persistence_scheduled_task.yml
title: Persistence via Scheduled Task Creation
id: f6a7b8c9-d0e1-2345-fabc-456789012345
status: stable
description: Detects creation of scheduled tasks, a common persistence mechanism
references:
  - https://attack.mitre.org/techniques/T1053/005/
author: Kousik Gunasekaran
date: 2024/01/01
tags:
  - attack.persistence
  - attack.t1053.005
  - attack.privilege_escalation
logsource:
  category: process_creation
  product: windows
detection:
  selection_schtasks:
    Image|endswith: '\schtasks.exe'
    CommandLine|contains: '/create'
  selection_at:
    Image|endswith: '\at.exe'
  selection_xml:
    Image|endswith: '\schtasks.exe'
    CommandLine|contains: '/xml'
  condition: selection_schtasks or selection_at or selection_xml
falsepositives:
  - Legitimate software installations creating scheduled tasks
  - IT automation and patch management tools
level: medium
